# launchd.google-calendar-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>launched.google-calendar-sync</string>
    <key>ProgramArguments</key>
    <array>
      <string>sh</string>
      <string>-c</string>
      <string>google-cal.sh</string>
    </array>
    <key>StandardErrorPath</key>
    <string>/Users/shane/tmp</string>
    <key>StandardOutPath</key>
    <string>/Users/shane/tmp</string>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>WorkingDirectory</key>
    <string>/Users/shane</string>
  </dict>
</plist>